# backend/maisha-api/config/systemd-maisha-queue.service
[Unit]
Description=Maisha Queue Worker
After=network.target

[Service]
Type=simple
User=dan
WorkingDirectory=/home/dan/Development/code/Wu-Tang/flask/January/maisha/backend/maisha-api
ExecStart=/usr/bin/php artisan queue:work database --queue=default --tries=3 --timeout=90 --max-jobs=1000 --max-time=3600 --sleep=3
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
